# deploy/systemd/bagipangan-scheduler.timer
[Unit]
Description=Run BagiPangan Laravel scheduler every minute

[Timer]
OnCalendar=*-*-* *:*:00
Persistent=true
Unit=bagipangan-scheduler.service

[Install]
WantedBy=timers.target
